# private_dot_config/scripts/ARCH/executable_pulsemods.service
[Unit]
Description=add remapped source to pulse

[Service]
ExecStart=/home/ralldi/.config/scripts/pulsemods.sh

[Install]
WantedBy=default.target
